# deploy/iconsai-scraping.service
[Unit]
Description=IconsAI Scraping API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/iconsai-scraping
Environment=PATH=/opt/iconsai-scraping/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/opt/iconsai-scraping/.env
ExecStart=/opt/iconsai-scraping/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
